# files/etc/hotplug.d/npt/25-restart-collectd
#!/bin/sh

# 1. Check if the service is enabled in /etc/config/ (returns 0 if enabled)
if /etc/init.d/collectd enabled; then
    # 2. Check if the process is actually running using pgrep.
    #    We check specifically for 'collectd' to ensure we don't start a stopped service.
    if pgrep -x "collectd" > /dev/null; then
        logger -t "collectd-fix" "Time sync event detected. Restarting running collectd instance to prevent RRD errors..."
        /etc/init.d/collectd restart
    else
        logger -t "collectd-fix" "Time sync event detected, but collectd is not running. Skipping restart."
    fi

else
    # Optional debug log
   logger -t "collectd-fix" "collectd is disabled. Skipping restart."
fi