# files/etc/uci-defaults/00-disable-netdata
#!/bin/sh

# ENSURE NETDATA IS DISABLED IN INIT SYSTEM
# We only want it starting via rc.local, never S50.
/etc/init.d/netdata disable



exit 0