# files/etc/uci-defaults/00-setup-resources
#!/bin/sh

# 1. FORCE KILL NETDATA
# The standard stop command fails under high load. We must use killall -9.
if pidof netdata >/dev/null; then
    killall -9 netdata
    /etc/init.d/netdata disable
fi

# 2. Fix Collectd Errors
mkdir -p /tmp/rrd

# 3. Setup Delayed Netdata Start in rc.local
grep -q "netdata" /etc/rc.local || {
    sed -i '/exit 0/d' /etc/rc.local
    echo "" >> /etc/rc.local
    echo "# Delayed Netdata start" >> /etc/rc.local
    echo "(sleep 60; /etc/init.d/netdata start) &" >> /etc/rc.local
    echo "exit 0" >> /etc/rc.local
}

exit 0