# files/etc/uci-defaults/00-setup-resources
#!/bin/sh

# 1. DISABLE NETDATA (Critical Fix)
# We must disable it regardless of whether it's running or not.
# This removes the S50 startup link so it won't start later in this boot.
if [ -f /etc/init.d/netdata ]; then
    /etc/init.d/netdata stop
    /etc/init.d/netdata disable
fi

# 2. Kill if running (Safety Net)
# In case we are in a loop or manual restart, kill any existing instance.
if pidof netdata >/dev/null; then
    killall -9 netdata
fi

# 3. Fix Collectd Errors
# Create the missing directory to stop the log spam
mkdir -p /tmp/rrd

exit 0