# files/etc/uci-defaults/z70-fix-unbound-scripts-execution
#!/bin/sh

# Ensure unbound scripts are executable
chmod +x /usr/lib/unbound/*

exit 0
