# files/etc/uci-defaults/z70-fix-default-shell
#!/bin/sh

# change root shell to bash
sed -i -E s/'^(root:.*:\/bin\/)ash'/'\1bash'/ /etc/passwd

exit 0
