# files/etc/uci-defaults/z8-plexmediaserver
#!/bin/sh

{ [ -f /etc/should_keep_old_uci_config ]; } || {

touch /etc/config/plexmediaserver

}
exit 0