# files/etc/uci-defaults/z98-custom
#!/bin/sh

# partial fix for ssl
grep -q 'net.netfilter.nf_conntrack_tcp_no_window_check' </etc/sysctl.conf || echo 'net.netfilter.nf_conntrack_tcp_no_window_check=1' >>/etc/sysctl.conf

# change root shell to bash
sed -i -E s/'^(root:.*:\/bin\/)ash'/'\1bash'/ /etc/passwd

# disable plexmediaserver starting.
if [ -f /etc/init.d/plexmediaserver ]; then
  /etc/init.d/plexmediaserver disable
  /etc/init.d/plexmediaserver stop
fi

# disable tailscale starting.
if [ -f /etc/init.d/tailscale ]; then
  /etc/init.d/tailscale disable
  /etc/init.d/tailscale stop
fi

# disable ddns starting.
if [ -f /etc/init.d/ddns ]; then
  /etc/init.d/ddns disable
  /etc/init.d/ddns stop
fi

# disable zerotier starting.
if [ -f /etc/init.d/zerotier ]; then
  /etc/init.d/zerotier disable
  /etc/init.d/zerotier stop
fi

# add avahi group
if [ -f /etc/init.d/avahi-daemon ]; then
  echo "avahi:x:105:avahi" >> /etc/group
  echo "avahi:*:105:105:avahi:/var/run/avahi-daemon:/bin/false" >> /etc/passwd
  /etc/init.d/avahi-daemon restart
fi

if [ -f /etc/init.d/ksmbd ]; then
   # setup SMBUSER for ksmbd
   grep -q 'SMBUSER' </etc/group || echo 'SMBUSER:x:2468:SMBUSER' >>/etc/group
   grep -q 'SMBUSER' </etc/shadow || echo 'SMBUSER:$5$CN2s9xwg8.EvN2Mp$ZZDnr1e/reUUai15L62aRCpLNRvts70dVwk.1UX2hB.:19818:0:99999:7:::' >>/etc/shadow
   grep -q 'SMBUSER' </etc/passwd || echo 'SMBUSER:*:2468:2468:SMBUSER:/mnt:/bin/false' >>/etc/passwd

   ksmbd.adduser -a SMBUSER -p SMBPASSWORD

   # Check if the config template is the same as default. If it is, change it for our default.
   { [ -f /etc/ksmbd/ksmbd.conf.template ] && sed s/'[[:space:]]'//g /etc/ksmbd/ksmbd.conf.template | grep -qE '.+' && [ -f /rom/etc/ksmbd/ksmbd.conf.template ] && sed s/'[[:space:]]'//g /rom/etc/ksmbd/ksmbd.conf.template | grep -qE '.+' && ! [ "$(cat /etc/ksmbd/ksmbd.conf.template)" = "$(cat /rom/etc/ksmbd/ksmbd.conf.template)" ]; } || {
      mv -f /etc/ksmbd/ksmbd.conf.template.example /etc/ksmbd/ksmbd.conf.template
   }
fi

mkdir -p /var/lock

# fix btop "--utf-force" bad argument if it exist
if [ -f /etc/profile.d/btop.sh ]; then
   sed -i 's/--utf-force/--force-utf/g' /etc/profile.d/btop.sh
fi

# generate stuff for unbound
if [ -f /etc/init.d/unbound ]; then
   mkdir -p /etc/unbound
   [ -f /etc/ssl/certs/ca-certificates.crt ] && ln -sf /etc/ssl/certs/ca-certificates.crt /etc/unbound/ca-certificates.crt

   mkdir -p /var/lib/unbound
   chown -R unbound:unbound /etc/unbound /var/lib/unbound
   
   # Only copy root.key if it exists (it might not if we skipped anchor)
   if [ -f /etc/unbound/root.key ]; then
       cp -p /etc/unbound/root.key /var/lib/unbound/
   fi

   unbound-control-setup -d /etc/unbound
   chown -R unbound:unbound /etc/unbound /var/lib/unbound
   uci set unbound.ub_main.unbound_control='2'
   uci commit unbound
   /etc/init.d/unbound restart
fi

# use chronyd
if [ -f /etc/init.d/chronyd ]; then
   /etc/init.d/sysntpd disable
   /etc/init.d/sysntpd stop
   
   # Disable Cron to avoid Busybox parsing bugs.
   # Chronyd handles periodic sync automatically after the initial hotplug trigger.
   if [ -f /etc/init.d/cron ]; then
       /etc/init.d/cron stop
       /etc/init.d/cron disable
   fi
fi

# fix insmod
if [ -f /etc/init.d/qca-nss-mirred ]; then
   # Improved sed command with proper quoting
   grep -qF 'act_nssmirred.ko' /etc/init.d/qca-nss-mirred || sed -i 's/act_nssmirred/act_nssmirred.ko/g' /etc/init.d/qca-nss-mirred
fi
if [ -f /sbin/kmodloader ]; then
   ln -sf /sbin/kmodloader /sbin/insmod
fi

# tweak ecm
if [ -f /etc/init.d/qca-nss-ecm ]; then
	uci set ecm.global.acceleration_engine='nss'
	uci set ecm.general.disable_gro='1'
	uci commit ecm
fi

# copy over ksmbd.conf.template if it doesn't exist
if [ ! -f /etc/ksmbd/ksmbd.conf.template ] && [ -f /etc/ksmbd/ksmbd.conf.template.example ]; then
    mv /etc/ksmbd/ksmbd.conf.template.example /etc/ksmbd/ksmbd.conf.template
fi

exit 0