# files/etc/uci-defaults/z99-done
#!/bin/sh

logger -p notice -t luci-defaults "First boot setup complete."

# --- MARK AS FIRST BOOT DONE ---
touch /etc/first_boot_done

exit 0