# files/etc/uci-defaults/z99-done
#!/bin/sh

# --- SAFETY CHECK (The Loop Breaker) ---
# If this file exists, we have already run. Exit immediately.
if [ -f /etc/first_boot_done ];
then
    exit 0
fi

# --- MARK AS DONE ---
touch /etc/first_boot_done

# --- LOGGING ---
logger -p notice -t luci-defaults "First boot setup complete. Applying changes live..."

# --- APPLY SETTINGS WITHOUT REBOOT ---

# 1. Reload the UCI configuration (notify procd)
/sbin/reload_config

# 2. Manually trigger SMP Affinity
# Normally this runs at S99 (boot), but since we aren't rebooting,
# we trigger it now to balance IRQs.
if [ -x /etc/init.d/smp_affinity ]; then
    /etc/init.d/smp_affinity start
    /etc/init.d/smp_affinity enable
fi

# 3. Restart/Reload Critical Services
# Network might flicker, but it won't be a full reboot
/etc/init.d/network reload
/etc/init.d/firewall restart
/etc/init.d/dnsmasq reload
/etc/init.d/odhcpd reload
/etc/init.d/cron restart
# Restart Unbound to apply z98-tweaks changes
if [ -x /etc/init.d/unbound ]; then
    /etc/init.d/unbound restart
fi

# 4. Restart NSS Subsystem (CRITICAL for z98-custom tweaks)
# z98-custom modifies qca-nss-ecm config and fixes qca-nss-mirred.
# We must restart them to apply the 'acceleration_engine=nss' settings.
if [ -x /etc/init.d/qca-nss-mirred ]; then
    /etc/init.d/qca-nss-mirred restart
fi

if [ -x /etc/init.d/qca-nss-ecm ]; then
    # Restart ECM to pick up the disable_gro and engine settings
    /etc/init.d/qca-nss-ecm restart
    # Flush current connections so they get picked up by the new engine
    echo 1 > /sys/kernel/debug/ecm/ecm_db/defunct_all
fi

if [ -x /etc/init.d/qca-nss-drv ]; then
    # Restart the NSS Driver so it picks up all changes from affinity
    /etc/init.d/qca-nss-drv restart
fi

# 5. Apply System Control tweaks (sysctl)
sysctl -p

# 6. Add flag for configuration already created
touch /etc/should_keep_old_uci_config

logger -p notice -t luci-defaults "Live configuration application complete."

exit 0