# files/etc/init.d/smp_affinity_split
#!/bin/sh /etc/rc.common

START=99

boot() {
    # [OPTIONAL] Load Balancing the secondary NSS Queue
    # We explicitly find the second nss_queue0 (IRQ 53) and move it to CPU1 (Mask 2)
    # This prevents CPU3 from doing ALL the work if load is high.

    # This finds the active nss_queue lines, sorts by interrupt count,
    # and takes the second busiest one to put on CPU1.
    # (Note: This is an advanced tweak. Only use if CPU3 is saturated.)
    logger -p notice -t smp_affinity "Applying Split NSS Affinity Strategy..."

    for irq in $(grep -i "nss_queue" /proc/interrupts | awk '{print $1}' | tr -d ':'); do
        # Check if the IRQ number is less than 50 (The first set, usually 43-46)
        if [ "$irq" -lt 50 ]; then
             echo "8" > "/proc/irq/$irq/smp_affinity"
             logger -p notice -t smp_affinity "  - NSS Queue IRQ $irq pinned to CPU3 (Mask 8)"
        else
             # IRQs 50+ (The second set, usually 53-56) go to CPU1
             echo "2" > "/proc/irq/$irq/smp_affinity"
             logger -p notice -t smp_affinity "  - NSS Queue IRQ $irq pinned to CPU1 (Mask 2)"
        fi
    done

    set_affinity 8 "nss_empty"

    logger -p notice -t smp_affinity "Shared NSS Queue with CPU1."
}
