# files/etc/rc.local
#!/bin/sh

# Create the RRD directory for statistics to prevent collectd errors on boot
mkdir -p /tmp/rrd

# Note: Time sync is now handled by /etc/hotplug.d/iface/20-chrony-sync
# This ensures boot completes immediately (Green Light) without waiting for internet.

exit 0