# files/etc/hotplug.d/iface/20-chrony-sync
#!/bin/sh

# Only run if the interface is 'wan' and the action is 'ifup'
[ "$INTERFACE" = "wan" ] || exit 0
[ "$ACTION" = "ifup" ] || exit 0

# Check if the tool exists and runs it
if [ -x /usr/sbin/chrony-hotplug ]; then
    logger -t chrony-hotplug "Interface $INTERFACE is up. Triggering time sync..."
    /usr/sbin/chrony-hotplug
fi