# files/etc/hotplug.d/iface/99-smp-vpn-affinity
#!/bin/sh

# 1. Exit immediately if not an "ifup" action
[ "$ACTION" = "ifup" ] || exit 0

# 2. Exit immediately if the master service is NOT enabled
#    This respects the user's decision to keep optimization off.
/etc/init.d/smp_affinity_vpn enabled || exit 0

MASK=8  # CPU3

# Helper to pin process by name (userspace or kernel thread)
pin_process() {
    local pattern="$1"
    local name="$2"

    # Find PIDs (pgrep -f matches full command line)
    for pid in $(pgrep -f "$pattern"); do
        if [ -d "/proc/$pid" ]; then
            taskset -p $MASK $pid >/dev/null 2>&1
            logger -p notice -t sm_vpn "Pinned $name (PID $pid) to CPU3"
        fi
    done
}

# --- Logic for specific VPNs ---

# A. WireGuard (Kernel Threads)
# Look for threads matching 'wg-crypt-[interface_name]'
if [ -n "$(pgrep -f "wg-crypt-$DEVICE")" ]; then
    pin_process "wg-crypt-$DEVICE" "WireGuard Thread"
fi

# B. Tailscale (Userspace)
# Tailscale benefits from CPU3 to share cache with NSS incoming packets
# Tailscale interfaces usually start with 'tailscale', like 'tailscale0'.
case "$DEVICE" in
    tailscale*)
        pin_process "tailscaled" "Tailscale Daemon"
        ;;
esac

# C. ZeroTier (Userspace)
# Like Tailscale, benefits from avoiding CPU0->CPU3 context switching
# ZeroTier interfaces usually start with 'zt', like 'zt87654321'
case "$DEVICE" in
    zt*)
        pin_process "zerotier-one" "ZeroTier Daemon"
        ;;
esac
